# yocto/recipes-immersive/immersive/files/immersive-net.service
[Unit]
Description=Immersive early network config (systemd-networkd from immersive.conf)
DefaultDependencies=no
After=local-fs.target
Before=network-pre.target systemd-networkd.service
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/immersive-net.sh

[Install]
WantedBy=sysinit.target
