# yocto/recipes-immersive/immersive/files/immersive-render.service
[Unit]
Description=Immersive 360 render node (GStreamer decode + GLES2 warp/blend -> KMS)
After=network-online.target immersive-role.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/immersive/render-node
# config.json is generated per-node by immersive-config.sh from immersive.conf
ExecStartPre=/usr/bin/immersive-config.sh
# release the framebuffer console so the GL stage can take DRM master for KMS:
# stop the tty getty and unbind the vt consoles from the framebuffer.
ExecStartPre=-/usr/bin/systemctl stop getty@tty1.service
ExecStartPre=-/bin/sh -c 'for vt in /sys/class/vtconsole/vtcon*; do echo 0 > "$vt/bind" 2>/dev/null; done'
# performance governor for steady 1080p60 decode
ExecStartPre=-/bin/sh -c 'echo performance > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor 2>/dev/null'
ExecStart=/usr/bin/python3 /opt/immersive/render-node/agent.py --config /run/immersive/config.json
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
